# linux/udev/50-keychron.rules
KERNEL=="hidraw5", SUBSYSTEM=="hidraw", ACTION=="add", RUN+="/home/thiago/bin/toggle_keychron_media_keys %p"
