# linux/udev/50-keychron.rules
KERNEL=="hidraw5", SUBSYSTEM=="hidraw", ACTION=="add", RUN+="/home/thiago/bin/fix_keychron_keyboard.sh %p"
