# linux/gnome-autostart/sudo.desktop
[Desktop Entry]
Type=Application
Exec=bash -c "sudo $HOME/bin/fix_keychron_keyboard.sh"
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[en_US]=Keychron Keyboard
Name=Keychron Keyboard
Comment[en_US]=Fix Keychron Keyboard
Comment=Fix Keychron Keyboard
